# roles/ghost/templates/ghost.container.j2
[Unit]
Description=Ghost Blog
After=network-online.target

[Container]
Image=docker.io/library/ghost:latest
AutoUpdate=registry

ContainerName=ghost

# Ports
PublishPort=127.0.0.1:{{ ghost_port }}:2368

# Volumes
Volume=ghost-data.volume:/var/lib/ghost/content

# --- Secrets and Environment ---
# Injects the secret as the specific environment variable Ghost expects
Secret=ghost_smtp_password,type=env,target=mail__options__auth__pass

# Loads the rest of the mail config from the protected file
EnvironmentFile=/etc/ghost/ghost.env

# Non-sensitive Environment Variables
Environment=DOMAIN={{ ghost_domain_name }}
Environment=url=https://{{ ghost_domain_name }}
Environment=database__client=sqlite3
Environment=database__connection__filename=content/data/ghost.db
Environment=database__useNullAsDefault=true

# Permissions
NoNewPrivileges=true
DropCapability=all
# Required for the entrypoint to fix volume permissions
AddCapability=CHOWN
AddCapability=FOWNER
AddCapability=DAC_OVERRIDE
# Required for the entrypoint to drop to the 'node' user
AddCapability=SETUID
AddCapability=SETGID
Tmpfs=/tmp

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target