# roles/syncthing/templates/syncthing.container.j2
[Unit]
Description=Syncthing Server
Requires=ts-syncthing.service
After=ts-syncthing.service
BindsTo=ts-syncthing.service
PartOf=ts-syncthing.service

[Container]
Image=docker.io/syncthing/syncthing:latest
AutoUpdate=registry

ContainerName=syncthing

# Attach to Tailscale's network namespace
Network=container:ts-syncthing
HostName=syncthing

# Environment Variables
Environment=PUID={{ syncthing_user_id }}
Environment=PGID={{ syncthing_user_group_id }}

# Volumes
Volume=syncthing-data.volume:/var/syncthing

# --- Security Lockdown (Aligned with Upstream systemd unit) ---
NoNewPrivileges=true
DropCapability=all
# Required to fix volume permissions on startup (AmbientCapabilities=CAP_CHOWN CAP_FOWNER in systemd)
AddCapability=CHOWN
AddCapability=FOWNER
AddCapability=DAC_OVERRIDE
# Required for su-exec to drop privileges to the PUID/PGID
AddCapability=SETUID
AddCapability=SETGID

# Equivalent to ProtectSystem=strict and PrivateTmp=true
ReadOnly=true
Tmpfs=/tmp

# Healthcheck using shell to properly execute the pipe command
HealthCmd=CMD-SHELL curl -fkLsS -m 2 127.0.0.1:8384/rest/noauth/health | grep -o --color=never OK || exit 1
HealthInterval=15s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target